# iac/playbooks/templates/webapp.service.j2
[Unit]
Description=DevOps WebApp
After=network.target redis-server.service
Wants=redis-server.service

[Service]
Type=simple
WorkingDirectory=/opt/webapp
Environment=PORT=3000
Environment=REDIS_HOST=127.0.0.1
Environment=REDIS_PORT=6379
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=3
User=vagrant

[Install]
WantedBy=multi-user.target
